# systemd/system/pia-pf-refresh@.timer
[Unit]
Description=Refresh PIA port forwarding timer on %I

[Timer]
OnCalendar=*-*-* *:00/15:00

[Install]
WantedBy=timers.target
